# deploy/systemd/minslab-monitor.timer
[Unit]
Description=Run MinsLab availability probe every minute

[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
AccuracySec=5s
Persistent=true
Unit=minslab-monitor.service

[Install]
WantedBy=timers.target
